# install/etc/profile.d/pythonenv.sh
source /opt/venv/djangoldap/bin/activate
source /opt/djangoldap/bin/setenv.sh
